# nc-ddns.timer
[Unit]
Description=Namecheap Dynamic DNS timer

[Timer]
OnCalendar=*:*

[Install]
WantedBy=timers.target
